# deploy/hoaproxy.service
# systemd unit that brings up the docker-compose stack on boot.
# Install with:
#   sudo cp deploy/hoaproxy.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hoaproxy
[Unit]
Description=hoaproxy.org docker-compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=hoaproxy
Group=hoaproxy
WorkingDirectory=/home/hoaproxy/hoaproxy/deploy
ExecStartPre=/usr/bin/docker compose -f docker-compose.prod.yml pull --ignore-pull-failures
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d --build
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
